# lib/systemd/system/xvnc.service
[Unit]
Description=x11vnc start
After=multi-user.target

[Service]
Type=simple
ExecStart=/bin/vnc-fun

[Install]
WantedBy=multi-user.target

